# dashboard-core/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:web]
# use a shell so ENV_ADMIN_PORT is read from the environment rather than supervisor interpolation
command=/bin/sh -c 'python3 -m http.server "$ENV_ADMIN_PORT" --directory /app/admin'
autostart=true
autorestart=true
stdout_logfile=/var/log/dashboard-web.log
stderr_logfile=/var/log/dashboard-web.err.log

[program:api]
command=python3 /app/admin_api.py
autostart=true
autorestart=true
stdout_logfile=/var/log/dashboard-api.log
stderr_logfile=/var/log/dashboard-api.err.log

[program:cache_fetcher]
command=python3 /app/cache_fetcher.py
autostart=true
autorestart=true
stdout_logfile=/var/log/dashboard-cache.log
stderr_logfile=/var/log/dashboard-cache.err.log

[program:dashboard_info]
command=python3 /app/dashboard_info.py
autostart=true
autorestart=true
stdout_logfile=/var/log/dashboard-info.log
stderr_logfile=/var/log/dashboard-info.err.log
